# systemd/cobien-launcher.service
[Unit]
Description=CoBien Launcher (main runtime)
After=network-online.target
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple
WorkingDirectory=%h/cobien/cobien-furniture-app-launcher
EnvironmentFile=-%h/cobien/cobien-furniture-app-launcher/cobien.env
Environment=COBIEN_LAUNCHER_ROOT=%h/cobien/cobien-furniture-app-launcher
ExecStart=/bin/bash -lc 'exec "$COBIEN_LAUNCHER_ROOT/cobien-launcher.sh" --mode launch --non-interactive --yes'
Restart=always
RestartSec=5
KillMode=mixed
TimeoutStopSec=20
