# systemd/cobien-update.service
[Unit]
Description=CoBien one-shot update check
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/cobien/cobien-furniture-app-launcher
EnvironmentFile=-%h/cobien/cobien-furniture-app-launcher/cobien.env
Environment=COBIEN_LAUNCHER_ROOT=%h/cobien/cobien-furniture-app-launcher
ExecStart=/bin/bash -lc 'exec "$COBIEN_LAUNCHER_ROOT/cobien-launcher.sh" --mode update-once --non-interactive --yes'
